# updater/vigilant-updater.service
[Unit]
Description=Vigilant Sensor - Remote Updater
Documentation=https://github.com/Vigilant-devs/vigilant-sensor-updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/vigilant/scripts/vigilantsensor
ExecStart=/vigilant/scripts/vigilantsensor/updater/vigilant-updater.sh
# Logging handled inside the script (compatible with systemd 239 / Rocky Linux 8)
StandardOutput=journal
StandardError=journal
